# .config/autostart/xcompmgr.desktop
[Desktop Entry]
Exec=/usr/bin/xcompmgr
Name=xcompmgr
Comment=sample X compositing manager
Type=Application
